# files/usr/lib/systemd/system/check-battery.service
[Unit]
Description=Check if mainboard battery is Ok
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/base-scripts/check-battery
