# files/etc/profile.d/malloc-debug.csh
# malloc debugging helpers (see man malloc), remove again for RC1 - meissner@suse.de
setenv MALLOC_CHECK_ 2
setenv MALLOC_PERTURB_ E
