# files/etc/profile.d/malloc-debug.sh
# malloc debugging helpers (see man malloc), remove again for RC1 - meissner@suse.de
export MALLOC_CHECK_=3
export MALLOC_PERTURB_=E
