# files/etc/profile.d/xdg-enviroment.csh
if ( ${?XDG_DATA_DIRS} ) then
    setenv XDG_DATA_DIRS /usr/share/:/usr/local/share/:/etc/opt/kde3/share/:/opt/kde3/share/:/opt/gnome/share/
endif
