# files/usr/etc/profile.d/malloc-debug.csh
# remove this file before shipment
setenv MALLOC_CHECK_ 3
setenv MALLOC_PERTURB_ 69
