# files/usr/etc/profile.d/xdg-environment.csh
# keep in sync with /usr/lib/environment.d/50-xdg.conf
setenv XDG_CONFIG_DIRS "/etc/xdg:/usr/local/etc/xdg:/usr/etc/xdg"
